# roles/setspyfm/templates/setspyfm-backend.service.j2
[Unit]
Description=SetSpyFM Backend Service
After=network.target

[Service]
KillMode=none
ExecStart=/home/{{ created_user }}/setspyfm-backend/setspy_api/venv/bin/python3.6 manage.py runserver
Restart=on-failure
WorkingDirectory=/home/{{ created_user }}/setspyfm-backendr/

[Install]
WantedBy=multi-user.target
